# config/.config/systemd/user/tldr-update.timer
[Unit]
Description=Run tldr update daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
